# tailfsrv.ini
[program:tailflog_ws_srv]
command=node tailfsrv.js
directory=/home/yuqing/tailflog
user=yuqing
redirect_stderr=true
stdout_logfile=/home/yuqing/tailflog/output.log
